# make_img_sample/kbox9_iso_build/rc.local
#!/bin/sh -e

echo "log" > /usr/local/log

cd /usr/new_file/
chmod 750 start.sh
./start.sh
exit 0
